# hello_cloud.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/hello_linux.py

[Install]
WantedBy=default.target